# .config/systemd/user/daily.service
[Unit]
Description=Daily routine

[Service]
ExecStart=/home/koneko/.local/bin/systemd-daily
