# .config/systemd/user/daily.timer
[Unit]
Description=Daily routine

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
Persistent=true
