# .config/systemd/user/weekly.service
[Unit]
Description=Weekly routine

[Service]
ExecStart=/home/koneko/.local/bin/systemd-weekly
